# Detection Conversion Files/ChatGPT-4o_Converted_Outputs Test 1 Temp 0.2/Failed/linux_at_application_execution.yml-output.json
{
  "actions": [],
  "author": ["OpenAI GPT4-o"],
  "description": "The following analytic detects the execution of the \"At\" application in Linux, which can be used by attackers to create persistence entries on a compromised host. This detection leverages data from Endpoint Detection and Response (EDR) agents, focusing on process names and parent process names associated with \"at\" or \"atd\". This activity is significant because the \"At\" application can be exploited to maintain unauthorized access or deliver additional malicious payloads. If confirmed malicious, this behavior could lead to data theft, ransomware attacks, or other severe consequences. Immediate investigation is required to determine the legitimacy of the execution and mitigate potential risks. Converted with OpenAI GPT4-o.",
  "enabled": true,
  "exceptions_list": [],
  "false_positives": ["Administrator or network operator can use this application for automation purposes. Please update the filter macros to remove false positives."],
  "from": "now-360s",
  "interval": "5m",
  "investigation_fields": {
    "field_names": []
  },
  "license": "",
  "max_signals": 100,
  "meta": {},
  "name": "OpenAI GPT4-o Linux At Application Execution",
  "note": "",
  "outcome": "exactMatch",
  "references": [
    "https://attack.mitre.org/techniques/T1053/001/",
    "https://www.linkedin.com/pulse/getting-attacker-ip-address-from-malicious-linux-job-craig-rowland/"
  ],
  "related_integrations": [],
  "required_fields": [],
  "risk_score": 50,
  "risk_score_mapping": [],
  "rule_id": "linux-at-application-execution",
  "rule_name_override": "",
  "severity": "medium",
  "severity_mapping": [],
  "tags": ["OpenAI GPT4-o", "Linux Privilege Escalation", "Linux Persistence Techniques", "Linux Living Off The Land", "Scheduled Tasks"],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053",
          "name": "Scheduled Task/Job",
          "reference": "https://attack.mitre.org/techniques/T1053/",
          "subtechnique": [
            {
              "id": "T1053.002",
              "name": "At",
              "reference": "https://attack.mitre.org/techniques/T1053/002/"
            }
          ]
        }
      ]
    }
  ],
  "throttle": "no_actions",
  "timeline_id": "",
  "timeline_title": "",
  "timestamp_override": "",
  "timestamp_override_fallback_disabled": true,
  "to": "now",
  "version": 1,
  "language": "eql",
  "query": "process where process.name in (\"at\", \"atd\") or process.parent.name in (\"at\", \"atd\")",
  "type": "eql",
  "alert_suppression": {
    "duration": {
      "unit": "s",
      "value": 3600
    },
    "group_by": ["host.name"],
    "missing_fields_strategy": "doNotSuppress"
  },
  "data_view_id": "",
  "event_category_override": "",
  "filters": [],
  "index": ["logs-*"],
  "tiebreaker_field": "",
  "timestamp_field": "@timestamp"
}