# Detection Conversion Files/Gemini-2.0-flash-001_Converted_Outputs Test 4 Temp 0.8/Failed/Linux_Edit_Cron_Table_Parameter.json
{
  "actions": [],
  "author": [
    "Google Gemini 2.0 Flash 001",
    "Teoderick Contreras, Splunk"
  ],
  "building_block_type": "default",
  "description": "Google Gemini 2.0 Flash 001: The following analytic detects the suspicious editing of cron jobs in Linux using the crontab command-line parameter (-e). It identifies this activity by monitoring command-line executions involving 'crontab' and the edit parameter. This behavior is significant for a SOC as cron job manipulations can indicate unauthorized persistence attempts or scheduled malicious actions. If confirmed malicious, this activity could lead to system compromise, unauthorized access, or broader network compromise.",
  "enabled": true,
  "false_positives": [
    "Administrator or network operator can use this application for automation purposes. Please update the filter macros to remove false positives."
  ],
  "from": "now-60m",
  "interval": "5m",
  "license": "Elastic License v2",
  "max_signals": 100,
  "name": "Google Gemini 2.0 Flash 001: Linux Edit Cron Table Parameter",
  "note": "This rule was converted with Google Gemini 2.0 Flash 001",
  "output_index": ".siem-signals-default",
  "references": [
    "https://attack.mitre.org/techniques/T1053/003/"
  ],
  "risk_score": 47,
  "rule_id": "linux_edit_cron_table_parameter",
  "severity": "medium",
  "tags": [
    "Google Gemini 2.0 Flash 001",
    "Linux Privilege Escalation",
    "Linux Persistence Techniques",
    "Linux Living Off The Land",
    "Scheduled Tasks",
    "Asset Type: Endpoint",
    "MITRE ATT&CK: T1053.003"
  ],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053",
          "name": "Scheduled Task/Job",
          "reference": "https://attack.mitre.org/techniques/T1053/",
          "subtechnique": [
            {
              "id": "T1053.003",
              "name": "Cron",
              "reference": "https://attack.mitre.org/techniques/T1053/003/"
            }
          ]
        }
      ]
    }
  ],
  "throttle": "no_actions",
  "type": "eql",
  "version": 1,
  "language": "eql",
  "query": """sequence by host.name
  [process where event.type == "start" and process.name == "crontab" and process.args in ("-e")]
  [process where event.type == "start" and process.name == "vi" and process.parent.name == "crontab"]""",
  "index": [
    "auditbeat-*",
    "logs-endpoint.events.*",
    "endgame-*",
    "logs-system.*"
  ]
}