# Detection Conversion Files/ChatGPT-4o_Converted_Outputs Test 1 Temp 0.2/Completed/linux_edit_cron_table_parameter.yml-output.json
{
  "name": "OpenAI GPT4-o Linux Edit Cron Table Parameter",
  "description": "The following analytic detects the suspicious editing of cron jobs in Linux using the crontab command-line parameter (-e). It identifies this activity by monitoring command-line executions involving 'crontab' and the edit parameter. This behavior is significant for a SOC as cron job manipulations can indicate unauthorized persistence attempts or scheduled malicious actions. If confirmed malicious, this activity could lead to system compromise, unauthorized access, or broader network compromise. Converted with OpenAI GPT4-o.",
  "author": ["Teoderick Contreras", "OpenAI GPT4-o"],
  "tags": ["OpenAI GPT4-o", "Linux Privilege Escalation", "Linux Persistence Techniques", "Linux Living Off The Land", "Scheduled Tasks"],
  "risk_score": 42,
  "severity": "medium",
  "type": "eql",
  "language": "eql",
  "query": "process where process.name == \"crontab\" and process.args : \"*-e*\"",
  "index": ["logs-*"],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053",
          "name": "Scheduled Task/Job",
          "reference": "https://attack.mitre.org/techniques/T1053/",
          "subtechnique": [
            {
              "id": "T1053.003",
              "name": "Cron",
              "reference": "https://attack.mitre.org/techniques/T1053/003/"
            }
          ]
        }
      ]
    }
  ],
  "references": ["https://attack.mitre.org/techniques/T1053/003/"],
  "enabled": true,
  "interval": "5m",
  "max_signals": 100,
  "version": 1,
  "exceptions_list": [],
  "actions": []
}