# Detection Conversion Files/ChatGPT-4o_Converted_Outputs Test 1 Temp 0.2/Completed/linux_add_files_in_known_crontab_directories.yml-output.json
{
  "name": "OpenAI GPT4-o Linux Add Files In Known Crontab Directories",
  "description": "The following analytic detects unauthorized file creation in known crontab directories on Unix-based systems. It leverages filesystem data to identify new files in directories such as /etc/cron* and /var/spool/cron/*. This activity is significant as it may indicate an attempt by threat actors or malware to establish persistence on a compromised host. If confirmed malicious, this could allow attackers to execute arbitrary code at scheduled intervals, potentially leading to further system compromise and unauthorized access to sensitive information. Converted with OpenAI GPT4-o.",
  "author": ["Teoderick Contreras", "OpenAI GPT4-o"],
  "tags": ["OpenAI GPT4-o"],
  "risk_score": 42,
  "severity": "medium",
  "type": "eql",
  "language": "eql",
  "query": "file where file.path : (\"*/etc/cron*\", \"*/var/spool/cron/*\")",
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053",
          "name": "Scheduled Task/Job",
          "reference": "https://attack.mitre.org/techniques/T1053/",
          "subtechnique": [
            {
              "id": "T1053.003",
              "name": "Cron",
              "reference": "https://attack.mitre.org/techniques/T1053/003/"
            }
          ]
        }
      ]
    }
  ],
  "references": [
    "https://www.sandflysecurity.com/blog/detecting-cronrat-malware-on-linux-instantly/",
    "https://www.cyberciti.biz/faq/how-do-i-add-jobs-to-cron-under-linux-or-unix-oses/"
  ],
  "index": ["filebeat-*"],
  "from": "now-6h",
  "interval": "5m",
  "enabled": true,
  "max_signals": 100,
  "version": 1
}