# Detection Conversion Files/Gemini-2.0-flash-001_Converted_Outputs Test 3 Temp 0.6 TopP 0.9/Failed/Linux_At_Application_Execution.json
{
  "actions": [],
  "alias_purpose": null,
  "alias_target_id": null,
  "author": [
    "Teoderick Contreras, Splunk",
    "Google Gemini 2.0 Flash 001"
  ],
  "building_block_type": null,
  "description": "The following analytic detects the execution of the \"At\" application in Linux, which can be used by attackers to create persistence entries on a compromised host. This detection leverages data from Endpoint Detection and Response (EDR) agents, focusing on process names and parent process names associated with \"at\" or \"atd\". This activity is significant because the \"At\" application can be exploited to maintain unauthorized access or deliver additional malicious payloads. If confirmed malicious, this behavior could lead to data theft, ransomware attacks, or other severe consequences. Immediate investigation is required to determine the legitimacy of the execution and mitigate potential risks. This rule was converted with Google Gemini 2.0 Flash 001.",
  "enabled": true,
  "exceptions_list": [],
  "false_positives": [
    "Administrator or network operator can use this application for automation purposes. Please update the filter macros to remove false positives."
  ],
  "from": "now-60m",
  "interval": "1m",
  "investigation_fields": null,
  "license": null,
  "max_signals": null,
  "meta": {
    "version": 6
  },
  "name": "Google Gemini 2.0 Flash 001 Linux At Application Execution",
  "namespace": "default",
  "note": null,
  "outcome": null,
  "output_index": ".siem-signals-default",
  "references": [
    "https://attack.mitre.org/techniques/T1053/001/",
    "https://www.linkedin.com/pulse/getting-attacker-ip-address-from-malicious-linux-job-craig-rowland/"
  ],
  "related_integrations": [],
  "required_fields": [
    {
      "name": "process.name",
      "type": "string"
    },
    {
      "name": "process.parent.name",
      "type": "string"
    }
  ],
  "response_actions": [],
  "risk_score": 42,
  "risk_score_mapping": [],
  "rule_id": "bf0a378e-5f3c-11ec-a6de-acde48001122",
  "rule_name_override": null,
  "setup": null,
  "severity": "low",
  "severity_mapping": [],
  "tags": [
    "Linux Privilege Escalation",
    "Linux Persistence Techniques",
    "Linux Living Off The Land",
    "Scheduled Tasks",
    "Endpoint",
    "T1053.002",
    "Google Gemini 2.0 Flash 001"
  ],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053",
          "name": "Scheduled Task/Job",
          "reference": "https://attack.mitre.org/techniques/T1053/",
          "subtechnique": [
            {
              "id": "T1053.001",
              "name": "At",
              "reference": "https://attack.mitre.org/techniques/T1053/001/"
            }
          ]
        }
      ]
    }
  ],
  "throttle": null,
  "timeline_id": null,
  "timeline_title": null,
  "timestamp_override": null,
  "timestamp_override_fallback_disabled": false,
  "to": "now",
  "version": 1,
  "language": "eql",
  "query": "sequence by host.id, user.id\n  [process where process.name in (\"at\", \"atd\") or process.parent.name in (\"at\", \"atd\")]",
  "type": "eql",
  "alert_suppression": {},
  "data_view_id": null,
  "event_category_override": null,
  "filters": [],
  "index": [
    "logs-endpoint.events.*"
  ],
  "tiebreaker_field": null,
  "timestamp_field": null
}