# Detection Conversion Files/Gemini-2.0-flash-001_Converted_Outputs Test 3 Temp 0.6 TopP 0.9/Failed/Linux_Possible_Cronjob_Modification_With_Editor.json
{
  "actions": [],
  "author": [
    "Google Gemini 2.0 Flash 001",
    "Teoderick Contreras"
  ],
  "description": "Google Gemini 2.0 Flash 001: The following analytic detects potential unauthorized modifications to\n  Linux cronjobs using text editors like \"nano,\" \"vi,\" or \"vim.\" It identifies this\n  activity by monitoring command-line executions that interact with cronjob configuration\n  paths. This behavior is significant for a SOC as it may indicate attempts at privilege\n  escalation or establishing persistent access. If confirmed malicious, the impact\n  could be severe, allowing attackers to execute damaging actions such as data theft,\n  system sabotage, or further network penetration. This rule was converted with Google Gemini 2.0 Flash 001.",
  "enabled": true,
  "false_positives": [
    "Administrator or network operator can use this commandline for automation purposes. Please update the filter macros to remove false positives."
  ],
  "from": "now-60m",
  "interval": "10m",
  "license": "Elastic License v2",
  "name": "Google Gemini 2.0 Flash 001: Linux Possible Cronjob Modification With Editor",
  "note": "The following analytic detects potential unauthorized modifications to\n  Linux cronjobs using text editors like \"nano,\" \"vi,\" or \"vim.\" It identifies this\n  activity by monitoring command-line executions that interact with cronjob configuration\n  paths. This behavior is significant for a SOC as it may indicate attempts at privilege\n  escalation or establishing persistent access. If confirmed malicious, the impact\n  could be severe, allowing attackers to execute damaging actions such as data theft, system sabotage, or further network penetration.",
  "references": [
    "https://attack.mitre.org/techniques/T1053/003/"
  ],
  "risk_score": 47,
  "rule_id": "dcc89bde-5f24-11ec-87ca-acde48001122",
  "severity": "medium",
  "tags": [
    "Google Gemini 2.0 Flash 001",
    "XorDDos",
    "Linux Living Off The Land",
    "Linux Privilege Escalation",
    "Scheduled Tasks",
    "Linux Persistence Techniques"
  ],
  "threat": [
    {
      "framework": "MITRE ATT&CK",
      "tactic": {
        "id": "TA0003",
        "name": "Persistence",
        "reference": "https://attack.mitre.org/tactics/TA0003/"
      },
      "technique": [
        {
          "id": "T1053.003",
          "name": "Scheduled Task/Job: Cron",
          "reference": "https://attack.mitre.org/techniques/T1053/003/"
        }
      ]
    }
  ],
  "type": "eql",
  "version": 1,
  "language": "eql",
  "query": """sequence by host.hostname, user.name
    [process where event.type == "start" and process.name in ("nano", "vi", "vim", "vim.basic") and process.args in ("/etc/cron*", "/var/spool/cron/*", "/etc/anacrontab*")]
    [process where event.type == "start" and process.name in ("nano", "vi", "vim", "vim.basic") and process.args in ("/etc/cron*", "/var/spool/cron/*", "/etc/anacrontab*")]""",
  "index": [
    "logs-endpoint.events.*"
  ],
    "data_view_id": "global"
}